# rel_b/etc/network/if-up.d/dns
#!/bin/sh

if [ $IFACE = "eth0" ]; then
	/etc/ssvconfig/init.d/network.sh final >/dev/null
fi
exit 0
