# rel_b/etc/network/if-up.d/ntpdate
#!/bin/sh

if [ $IFACE = "eth0" ]; then
    /etc/ssvconfig/init.d/time.sh settime >/dev/null
fi
exit 0
